# debian/postinst
#!/bin/bash
set -e
#DEBHELPER#
. /usr/share/debhelper/dh_installdeb
dh_systemd_enable --user
systemctl --user daemon-reload
systemctl --user enable dad_joke.timer --now
systemctl --user start dad_joke.service || true
exit 0
